# guardiankey-ssh/etc/cron.d/guardiankey
## GuardianKey - Unlock benned IPs
*/10 * * * * /usr/bin/python /usr/lib/guardiankey/gkUnlock.py
