# rpmbuild/BUILD/guardiankey-ssh-1/etc/cron.d/guardiankey
## GuardianKey - Unlock benned IPs
* * * * * /usr/bin/python /usr/lib/guardiankey/gkUnlock.py
